# files/supervisor-uwsgi.ini
[program:uwsgi]
user=uwsgi
group=uwsgi
stopsignal=INT
directory=/app
redirect_stderr=true
command=/usr/bin/uwsgi --ini /configs/uwsgi.ini
